# root/etc/cron.d/nextcloud
*/5 * * * *   apache /usr/bin/scl enable rh-php73 -- php -d memory_limit=512M -d apc.enable_cli=1 -f /usr/share/nextcloud/cron.php
