# root/etc/cron.d/nextcloud
*/15 * * * *   apache /opt/rh/rh-php71/root/usr/bin/php -f /usr/share/nextcloud/cron.php
